# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:dugong]
user=jiaxin
directory=/data/www/dugong
command=/home/jiaxin/.virtualenvs/dugong/bin/gunicorn -c /etc/default/gunicorn dugong.wsgi