# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:dugong]
directory=/data/www/dugong
command=/usr/local/bin/gunicorn -c /etc/default/gunicorn dugong.wsgi